# deploy/systemd/tasseograph-agent.service
[Unit]
Description=Tasseograph dmesg Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tasseograph agent --config /etc/tasseograph/agent.yaml
Restart=always
RestartSec=10
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tasseograph
EnvironmentFile=-/etc/tasseograph/agent.env

[Install]
WantedBy=multi-user.target
